# scripts/com.liveassistant.cursor-monitor.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.liveassistant.cursor-monitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/Desktop/Projects/LiveAssistant/scripts/cursor-daemon.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <!-- CHANGE THIS: Update to your actual project path -->
    <string>/Users/YOUR_USERNAME/Desktop/Projects/LiveAssistant/logs/cursor-daemon.log</string>
    
    <key>StandardErrorPath</key>
    <!-- CHANGE THIS: Update to your actual project path -->
    <string>/Users/YOUR_USERNAME/Desktop/Projects/LiveAssistant/logs/cursor-daemon.error.log</string>
    
    <key>WorkingDirectory</key>
    <!-- CHANGE THIS: Update to your actual project path -->
    <string>/Users/YOUR_USERNAME/Desktop/Projects/LiveAssistant</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>

